# ansible/roles/docker/templates/docker-prune.timer.j2
[Unit]
Description=Run Docker/Containerd Prune Weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
